# docker/feedback-by-page/crontab
@reboot cd /app && $(which python3) ./fetch_table.py > /proc/1/fd/1 2>/proc/1/fd/2
0 14 * * * cd /app && $(which python3) ./fetch_table.py > /proc/1/fd/1 2>/proc/1/fd/2
