# docker/train/crontab
@reboot cd /app && $(which python3) ./process_feedback.py > /proc/1/fd/1 2>/proc/1/fd/2
0 15 * * * cd /app && $(which python3) ./process_feedback.py > /proc/1/fd/1 2>/proc/1/fd/2
